# dist/p2.service
[Unit]
Description=Conversations Push Proxy

[Service]
User=p2
Group=p2
ExecStart=/usr/bin/java -jar /opt/p2-0.3.jar -c /etc/p2/config.json

[Install]
WantedBy=multi-user.target
